# process_list/process_list.service
[Unit]
Description=Evacuate engine service. Owner Alexander Semenov
After=network.target
 
[Service]
User=root
Group=root
Type=simple
Restart=always
ExecStart=/usr/local/bin/process_list
 
[Install]
WantedBy=multi-user.target
